# etc/cron.d/chk_condump_crond
# Run check_condump.sh every 30mins
*/30 * * * * root /sbin/afail.sh -k > /dev/null 2>&1
