# etc/cron.d/raid_crond
# Run periodic_raidcheck.py every 5mins as a fallback for sw_raidcheck.py
*/5 * * * * root /sbin/periodic_raidcheck.sh > /dev/null 2>&1

