# etc/cron.d/rcu_crond
# clean out /var/log/rcu/backup every other day
0 0 */2 * * root /sbin/rcu_backup_clean.sh > /dev/null 2>&1
