# etc/cron.d/samba_crond
# clean out /var/log/samba/cores every other day
0 0 */2 * * root /sbin/samba_core_clean.sh > /dev/null 2>&1
