# etc/cron.d/webasd_watchdog_crond
# $URL: svn://svn.nbttech.com/mgmt-fwk/branches/kauai_373_fix_branch/framework/src/base_os/common/image_files/webasd_watchdog_crond $
#
# check every minute
* * * * * root /sbin/webasd_watchdog.sh > /dev/null 2>&1
