# etc/profile.d/rvbd_profile.sh

export PATH=/usr/local/bin:/opt/rbt/bin:/opt/tms/bin:$PATH

if [ -r /opt/tms/release/build_version.sh ] ; then
   . /opt/tms/release/build_version.sh
fi

