# etc/init/serial-ttyS0.conf
# This service maintains a getty on /dev/ttyS0.

start on stopped rc RUNLEVEL=[2345]
stop on starting runlevel [016]

respawn
instance /dev/ttyS0
exec /sbin/con_getty -L ttyS0 9600 vt102
