# sdd/sdd6/etc/cron.d/ipmitool_crond
# sync ipmitool sel time with system time
0 0 * * * root /sbin/ipmitool_sync.sh > /dev/null 2>&1
