# sdd/sdd6/etc/cron.d/smb_crond
# backup smb data
0 0 * * * root /sbin/smb_backup.sh > /dev/null 2>&1
